# bsp/mach-tievm/apps/base/etc/inittab
::sysinit:/etc/rcS
::restart:/sbin/init
ttyO0::askfirst:-/bin/sh
